# launchd/com.contextd.nightly-digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.contextd.nightly-digest</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/amit/contextd/scripts/nightly-digest.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/amit/contextd/scripts</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>23</integer>
        <key>Minute</key>
        <integer>55</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/amit/Library/Logs/contextd-nightly-digest.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/amit/Library/Logs/contextd-nightly-digest.err</string>
</dict>
</plist>
