# launchd/com.contextd.consolidate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.contextd.consolidate</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/amit/contextd/scripts/consolidate-vault.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/amit/contextd/scripts</string>

    <key>StartInterval</key>
    <integer>21600</integer>

    <key>StandardOutPath</key>
    <string>/Users/amit/Library/Logs/contextd-consolidate.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/amit/Library/Logs/contextd-consolidate.log</string>
</dict>
</plist>
